# docker/workers/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:worker-virus-scan]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan docintern:consume virus-scan --sleep=3 --tries=3 --timeout=60
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=docintern
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/worker-virus-scan.log
stopwaitsecs=30

[program:worker-ocr-extraction]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan docintern:consume ocr-extraction --sleep=3 --tries=3 --timeout=300
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=docintern
numprocs=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/worker-ocr-extraction.log
stopwaitsecs=60

[program:worker-classification]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan docintern:consume classification --sleep=3 --tries=3 --timeout=120
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=docintern
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/worker-classification.log
stopwaitsecs=30

[program:worker-audit-log]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan docintern:consume audit-log --sleep=3 --tries=3 --timeout=30
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=docintern
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/worker-audit-log.log
stopwaitsecs=15

[program:worker-dlq]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan docintern:consume dead-letters --sleep=10 --tries=1 --timeout=60
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=docintern
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/worker-dlq.log
stopwaitsecs=15
